# docker-files.d/supervisord-celery-beat.conf
[supervisord]
nodaemon=true

[program:osp-celery-beat]
directory=/app
user=www-data
group=www-data
command=/usr/local/bin/celery -A app.celery beat
